# misc/phpplatesender.service.in
# phpplatesender systemd unit file
# This file is intended to work with debian distributions

[Unit]
Description=PHP based license plate notifications sender
After=network.target memcached.service openalpr-daemon.service beanstalkd.service
Requires=memcached.service openalpr-daemon.service beanstalkd.service
Wants=

[Service]
User=@RUN_USER@
Type=forking
ExecStart=/usr/bin/php @DATA_DIR@/daemon.php start
ExecStop=/usr/bin/php @DATA_DIR@/daemon.php stop
PIDFile=@RUN_DIR@/phpplatesender.pid

[Install]
WantedBy=multi-user.target
